# contrib/init/coinbitd.conf
description "Coinbit Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env COINBITD_BIN="/usr/bin/coinbitd"
env COINBITD_USER="coinbit"
env COINBITD_GROUP="coinbit"
env COINBITD_PIDDIR="/var/run/coinbitd"
# upstart can't handle variables constructed with other variables
env COINBITD_PIDFILE="/var/run/coinbitd/coinbitd.pid"
env COINBITD_CONFIGFILE="/etc/coinbit/coinbit.conf"
env COINBITD_DATADIR="/var/lib/coinbitd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # coinbitd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$COINBITD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run coinbitd."
        echo "The setting must appear in $COINBITD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Coinbit Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$COINBITD_PIDDIR"
    chmod 0755 "$COINBITD_PIDDIR"
    chown $COINBITD_USER:$COINBITD_GROUP "$COINBITD_PIDDIR"
    chown $COINBITD_USER:$COINBITD_GROUP "$COINBITD_CONFIGFILE"
    chmod 0660 "$COINBITD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$COINBITD_PIDFILE" \
    --chuid $COINBITD_USER:$COINBITD_GROUP \
    --exec "$COINBITD_BIN" \
    -- \
    -pid="$COINBITD_PIDFILE" \
    -conf="$COINBITD_CONFIGFILE" \
    -datadir="$COINBITD_DATADIR" \
    -disablewallet \
    -daemon

